# prusaslicer-novnc/rootfs/etc/services.d/03-prusaslicer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs PrusaSlicer
# ==============================================================================
bashio::log.info "Starting PrusaSlicer..."
# cd /opt || bashio::exit "Could not change dir."
export DISPLAY=:0
exec /bin/bash -c '/slic3r/slic3r-dist/prusa-slicer --datadir /slicer-configs/.config/PrusaSlicer/'
