# prusaslicer-novnc/rootfs/etc/services.d/02-openbox/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Openbox
# ==============================================================================
bashio::log.info "Starting Openbox..."
# cd /opt || bashio::exit "Could not change dir."
export DISPLAY=:0
exec /usr/bin/openbox
