# prusaslicer-novnc/rootfs/etc/services.d/02-openbox/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: PrusaSlicer NoVNC
# Runs Openbox
# ==============================================================================
bashio::net.wait_for 4325
bashio::log.info "Starting Openbox..."

export DISPLAY=:0
exec /usr/bin/openbox
