# prusaslicer-novnc/rootfs/etc/services.d/00-x11/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: PrusaSlicer NoVNC
# Runs X11
# ==============================================================================
bashio::log.info "Starting X11 server..."
# cd /opt || bashio::exit "Could not change dir."
exec /usr/bin/Xtigervnc -desktop "Prusaslicer" -localhost -rfbport 5900 -SecurityTypes None -AlwaysShared -AcceptKeyEvents -AcceptPointerEvents -AcceptSetDesktopSize -SendCutText -AcceptCutText :0
